# docker/supervisord.conf
[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=10
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
